# boot/gpu.service
[Unit]
Description=gpu test

[Service]
Type=simple
ExecStart=/usr/local/sbin/gpu_fps_test.run
